# systemd/discovery-server.service
[Unit]
Description=Fast DDS Discovery Server
After=network-online.target
Wants=network-online.target

[Service]
User=jetson1
Group=jetson1

# Run discovery server directly - no need to source full shell config
# Server ID 0 (-i 0) must match the prefix in super_client_template.xml
ExecStart=/bin/bash -c "source /opt/ros/humble/setup.bash && exec fastdds discovery -i 0 -p 11811"

# Always restart - discovery server is critical for ROS2 node communication
Restart=always
RestartSec=3

# Standard output goes to journald
StandardOutput=journal
StandardError=journal

# Give it time to bind to port
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target 